# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:opacweb]
command=/usr/local/bin/opacweb web

[program:memcached]
command=/usr/bin/memcached -p 11211 -u memcached -m 256 -c 1024 -t 4
autorestart=true

